# apps/linode-marketplace-prometheus-grafana/roles/loki/templates/loki.service.j2
[Unit]
Description=Grafana Loki log aggregation and storage service
After=network.target

[Service]
User=loki
Group=loki
Type=simple
ExecStart=/usr/local/bin/loki -config.file /etc/loki/loki-config.yml

[Install]
WantedBy=multi-user.target
