# apps/linode-marketplace-owncast/roles/owncast/templates/owncast.service.j2
[Unit]
Description=Owncast Live Streaming Server
After=network.target

[Service]
Type=simple
ExecStart=/opt/owncast/owncast
WorkingDirectory=/opt/owncast
User={{ username }}
Group={{ username }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 